# linux/kali/kali-etc-configs-master/sv/pcscd/run
#!/usr/bin/env /lib/runit/invoke-run

[ -e /etc/runit/verbose ] && echo "invoke-run: starting ${PWD##*/}"

exec /usr/sbin/pcscd --foreground $PCSCD_ARGS 2>&1
