# linux/kali/kali-etc-configs-master/xdg/autostart/xcape-super-key-bind.desktop
[Desktop Entry]
Type=Application
Name=Xcape Super Key Bind
Exec=xcape -e 'Super_L=Control_L|Escape'
Terminal=false
OnlyShowIn=XFCE;
